# system/10-add-comm1.rules
ACTION=="add", SUBSYSTEM=="input", RUN+="/root/ansible_add.sh %E{NAME}"
